# src/arm/j17084truckduck.conf
start on truck_duck_initialized
stop on (stop_j1708 or runlevel [016])

env LD_LIBRARY_PATH=/usr/local/lib

script
    chdir /lib/firmware/
    python2.7 /usr/local/bin/j17084truckduck_host.py
end script

